# com.tradetron.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradetron.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-u</string>
        <string>/Users/prasanjitdatta/Desktop/antigravity/tradetron_dashboard/dashboard_server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/prasanjitdatta/Desktop/antigravity/tradetron_dashboard</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/prasanjitdatta/Desktop/antigravity/tradetron_dashboard/dashboard.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/prasanjitdatta/Desktop/antigravity/tradetron_dashboard/dashboard_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
